# crontab
* * * * * /usr/bin/training_and_deploy.sh > /proc/1/fd/1 2>/proc/1/fd/2
# required empty line